# spacebrew.conf
#!upstart
description "spacebrew"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production
env PATH=/root/local/bin

exec node /var/www/spacebrew/node_server.js >> /var/log/spacebrew.log 2>&1
